# wireguard-clients/rootfs/etc/services.d/status/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard Clients
# Shows current WireGuard status
# ==============================================================================
sleep 30
bashio::log.info "Requesting current status from WireGuard Clients..."
exec wg show
